# systemd/energy-app.service
[Unit]
Description=Energy App Backend Service
After=network.target

[Service]
User=root
WorkingDirectory=/opt/energy-app
ExecStart=/usr/bin/python3 /opt/energy-app/server.py
Restart=always

[Install]
WantedBy=multi-user.target
